# back/utils/supervisor/dev-stat.conf
[program:dev-stat]
directory=/opt/dev-stat
command=/opt/dev-stat/start-app.sh
stderr_logfile = /var/log/dev-stat-stderr.log
stdout_logfile = /var/log/dev-stat-stdout.log
